# chip-configs/etc/systemd/system/plymouth-quit.service.d/retain-splash.conf
[Service]
ExecStart=
ExecStart=/bin/plymouth quit --retain-splash
